# nixie.service
[Service]
Restart=always
RestartSec=3
ExecStart=/usr/local/bin/nixie-clock

[Install]
WantedBy=multi-user.target

